# crontab
*/5 * * * * /home/arnaudchn/PROJET1/main_projet.py
